# package/luci-app-sheepfold-family-internet-control/root/etc/init.d/sheepfold
#!/bin/sh /etc/rc.common

START=90
USE_PROCD=1

start_service() {
        procd_open_instance
        procd_set_param command /usr/libexec/sheepfold/sheepfold-service
        procd_set_param respawn
        procd_close_instance
}

reload_service() {
        /usr/libexec/sheepfold/sheepfold-service reload
}

stop_service() {
        [ -x /usr/libexec/sheepfold/sheepfold-router-control ] &&
                /usr/libexec/sheepfold/sheepfold-router-control blocked-page-stop >/dev/null 2>&1 || true
}
